# etc/systemd/system/dnf-automatic.timer.d/schedule.conf
[Timer]
OnBootSec=
OnUnitActiveSec=
OnCalendar=*-*-* 03:00:00
AccuracySec=10m
RandomizedDelaySec=0
